# sysycfg/udev/rules.d/10-local.rules
KERNEL=="wifi0", RUN+="/sbin/ath_leds.sh"
